# systemd/domain-services.service
[Unit]
Description=domain-services — Postgres + MinIO via docker compose
Requires=docker.service
After=docker.service network-online.target
Wants=network-online.target

[Service]
Type=oneshot
RemainAfterExit=yes
WorkingDirectory=/opt/services
EnvironmentFile=/opt/services/.env
ExecStart=/usr/bin/make up
ExecStop=/usr/bin/make down
ExecReload=/usr/bin/make restart
TimeoutStartSec=300
TimeoutStopSec=120

[Install]
WantedBy=multi-user.target
