# systemd/domain-services-backup.timer
[Unit]
Description=domain-services backup diario 02:00 UTC

[Timer]
OnCalendar=*-*-* 02:00:00
Persistent=true
RandomizedDelaySec=300
Unit=domain-services-backup.service

[Install]
WantedBy=timers.target
