# services/systemd/domain-services-backup.service
[Unit]
Description=domain-services — backup diario
After=docker.service domain-services.service
Requires=domain-services.service

[Service]
Type=oneshot
WorkingDirectory=/opt/services
EnvironmentFile=/opt/services/.env
ExecStart=/opt/services/scripts/backup.sh
TimeoutStartSec=3600
Nice=10
IOSchedulingClass=best-effort
IOSchedulingPriority=7

[Install]
WantedBy=multi-user.target
