# services/systemd/domain-services.service
[Unit]
Description=domain-services — orquestación de 5 servicios via docker compose
Documentation=https://github.com/nunezlagos/domain/tree/services
Requires=docker.service
After=docker.service network-online.target
Wants=network-online.target

[Service]
Type=oneshot
RemainAfterExit=yes
WorkingDirectory=/opt/services/services
EnvironmentFile=/opt/services/.env
ExecStart=/usr/bin/make -C /opt/services/services up
ExecStop=/usr/bin/make -C /opt/services/services down
ExecReload=/usr/bin/make -C /opt/services/services restart SVC=all
TimeoutStartSec=300
TimeoutStopSec=120

[Install]
WantedBy=multi-user.target
